# .aws/assetmanagement.service
[Unit]
Description=Asset Management System - .NET Web Application
After=network.target
StartLimitIntervalSec=0

[Service]
Type=notify
# Run as ec2-user for security
User=ec2-user
Group=ec2-user

# Application settings
WorkingDirectory=/var/www/assetmanagement/current
ExecStart=/usr/bin/dotnet AssetManagementSystem.PL.dll
Restart=always
RestartSec=1

# Environment variables
Environment=ASPNETCORE_ENVIRONMENT=Production
Environment=ASPNETCORE_URLS=http://localhost:5000
Environment=DOTNET_PRINT_TELEMETRY_MESSAGE=false
Environment=DOTNET_CLI_TELEMETRY_OPTOUT=1

# Resource limits
MemoryMax=1G
TasksMax=100

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/assetmanagement/current/wwwroot/files
ReadWritePaths=/var/www/assetmanagement/current/Logs

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=assetmanagement

# Process management
KillMode=mixed
TimeoutStartSec=30
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
