# .aws/assetmanagement.service
[Unit]
Description=Asset Management System - ASP.NET Core Application
After=network.target
StartLimitIntervalSec=0

[Service]
Type=notify
User=ec2-user
Group=ec2-user
WorkingDirectory=/var/www/assetmanagement/current
ExecStart=/usr/local/bin/dotnet AssetManagementSystem.PL.dll
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5
TimeoutStopSec=90
SyslogIdentifier=assetmanagement
Environment=ASPNETCORE_ENVIRONMENT=Production
Environment=DOTNET_ROOT=/home/ec2-user/.dotnet
Environment=PATH=/home/ec2-user/.dotnet:/usr/local/bin:/usr/bin:/bin
Environment=ASPNETCORE_URLS=http://0.0.0.0:5000
Environment=DOTNET_PRINT_TELEMETRY_MESSAGE=false
Environment=DOTNET_CLI_TELEMETRY_OPTOUT=1

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/var/www/assetmanagement
ReadWritePaths=/tmp
ReadWritePaths=/var/log

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
